# templates/zuul.service.erb
# WARNING THIS FILE IS MANAGED BY PUPPET
[Unit]
Description=Zuul
After=network.target

[Service]
SyslogIdentifier=zuul
User=<%= @zuul_user %>
Group=<%= @zuul_user %>
Environment=VIRTUAL_ENV=<%= @venv_path %>
Environment=PATH=$VIRTUAL_ENV/bin:$PATH
ExecStart=<%= @venv_path %>/bin/zuul-server -d
KillSignal=SIGUSR1
ExecReload=/bin/kill -SIGHUP $MAINPID

[Install]
WantedBy=multi-user.target
